# deploy/local/aeon-aspects/launchagents/org.aeon.buzz-acp.voxis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>org.aeon.buzz-acp.voxis</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/AEON/Projects/buzz/target/release/buzz-acp</string>
    <string>--relay-url</string>
    <string>ws://localhost:3000</string>
    <string>--private-key-file</string>
    <string>/Volumes/AEON/Projects/buzz-data/keys/voxis.sk</string>
    <string>--expected-public-key</string>
    <string>e2b124c4728989a09d0f7df96cd2b1823635cc220ed9cdd4f97eaedc12a17fe4</string>
    <string>--agent-owner</string>
    <string>73ac8798fd9cedcc5d24645d7ed49332d94a28f2c937f4c1b638f92bf6e8e91f</string>
    <string>--agent-command</string>
    <string>/REQUIRES_FLEET/immutable-openclaw/bin/openclaw</string>
    <string>--agent-args</string>
    <string>acp,--session,agent:voxis:buzz-private,--require-existing,--token-file,/REQUIRES_FLEET/owned-token-file,--url,ws://127.0.0.1:18806,--provenance,meta+receipt,--no-prefix-cwd</string>
    <string>--agents</string>
    <string>1</string>
    <string>--subscribe</string>
    <string>config</string>
    <string>--config</string>
    <string>/Volumes/AEON/Projects/buzz/deploy/local/aeon-aspects/config/voxis.toml</string>
    <string>--respond-to</string>
    <string>owner-only</string>
    <string>--allowed-respond-to</string>
    <string>owner-only</string>
    <string>--no-memory</string>
    <string>--no-base-prompt</string>
    <string>--dedup</string>
    <string>queue</string>
    <string>--multiple-event-handling</string>
    <string>queue</string>
    <string>--relay-observer</string>
    <string>--trusted-inbound-envelope</string>
    <string>--permission-mode</string>
    <string>bypass-permissions</string>
    <string>--heartbeat-interval</string>
    <string>0</string>
    <string>--turn-liveness-secs</string>
    <string>10</string>
    <string>--idle-timeout</string>
    <string>900</string>
    <string>--max-turn-duration</string>
    <string>7200</string>
    <string>--context-message-limit</string>
    <string>12</string>
    <string>--max-turns-per-session</string>
    <string>0</string>
    <string>--turn-receipts</string>
    <string>--expected-gateway-session-key</string>
    <string>agent:voxis:buzz-private</string>
  </array>
  <key>WorkingDirectory</key><string>/Volumes/AEON/Projects/buzz</string>
  <key>RunAtLoad</key><false/>
  <key>KeepAlive</key><false/>
  <key>ProcessType</key><string>Background</string>
  <key>StandardOutPath</key><string>/Volumes/AEON/Projects/buzz-data/logs/voxis.buzz-acp.log</string>
  <key>StandardErrorPath</key><string>/Volumes/AEON/Projects/buzz-data/logs/voxis.buzz-acp.err.log</string>
</dict>
</plist>
